# srcs/ftps/supervisord.conf
[supervisord]
nodaemon=true
autorestart=true
autostart=true

[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd/vsftpd.conf
autorestart=true
autostart=true

[program:telegraf]
command=/usr/bin/telegraf -config /etc/telegraf.conf -config-directory /etc/tel
autorestart=true
autostart=true